# .config/systemd/user/_echo@.service
[Service]
Type=oneshot
ExecStart=/usr/bin/echo '%I'

[Install]
WantedBy=default.target
