# .config/systemd/user/dropbox.service
[Unit]
Description=Dropbox autostart

[Service]
Type=simple
ExecStart=/usr/local/bin/dropbox

[Install]
WantedBy=default.target

