# .config/systemd/user/environment.target
[Unit]
Description=Set environment variables
Requires=dbus.socket
IgnoreOnIsolate=true

