# .config/systemd/user/tmp-ammon@.service
[Service]
Type=oneshot
ExecStart=/usr/bin/mkdir -p /tmp/ammon/%i

[Install]
WantedBy=default.target
