# .config/systemd/user/mpd.service.d/override.conf
[Unit]
Requires=tmp-emmie@cache.service
